# ktrl.service
[Unit]
Description=ktrl

[Service]
ExecStart=/home/<your-user>/.cargo/bin/ktrl -d <keyboard-dev-path>
  
[Install]
WantedBy=multi-user.target
